# scripts-tmp/crontab
*/10 * * * * root /home/ec2-user/bin/idle-mc-shutdown
